# deploy/systemd/golewood-cron.service
[Unit]
Description=Golewood cron (calendar sync + promotion index)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/opt/golewood/.env
ExecStart=/opt/golewood/deploy/cron-hit.sh
